# stageastroplant0/03-kit-startup/files/pigpiod.service
[Unit]
Description=PiGPIOd Service

[Service]
ExecStart=pigpiod
Type=forking
StandardOutput=inherit
StandardError=inherit
User=root
